# gogs/rootfs/etc/s6/gogs/run
#!/bin/sh

su-exec app /app/gogs web -c /data/conf/app.ini
