# s6/maildev/run
#!/bin/sh

exec s6-setuidgid app maildev
