# s6/postgres/run
#!/bin/sh

set -e

./setup

su-exec app postgres
