# etc/s6/registry/run
#!/bin/sh

su-exec app docker-registry serve /etc/docker-registry/config.yml
