# s6/jenkins/run
#!/bin/sh

exec s6-setuidgid app java -Djava.awt.headless=true -jar /app/jenkins.war --webroot=/app/cache --httpPort=8080
