# rootfs/etc/s6/gogs/run
#!/bin/sh

exec s6-setuidgid app /app/gogs web -c /data/conf/app.ini
